# deploy/netrunner.service
[Unit]
Description=Lfmt.NetRunner deployment service
After=network.target

[Service]
WorkingDirectory=/opt/netrunner
ExecStart=/usr/bin/dotnet /opt/netrunner/Lfmt.NetRunner.dll
Restart=always
RestartSec=5
SyslogIdentifier=netrunner

User=netrunner
Group=netrunner

Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_NOLOGO=1

# NetRunner needs sudo for user/service management, so strict
# filesystem protection is not applicable here. Managed apps
# get full hardening in their own service files.
PrivateTmp=true
PrivateDevices=true

[Install]
WantedBy=multi-user.target
